# debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Create data directory structure
        # The binary uses relative paths (data/state.json, config/gateway.json)
        # WorkingDirectory=/var/lib/vpn-gateway in systemd unit
        mkdir -p /var/lib/vpn-gateway/data/wg
        mkdir -p /var/lib/vpn-gateway/data/unbound

        # Symlink config so the binary finds config/gateway.json
        if [ ! -e /var/lib/vpn-gateway/config ]; then
            ln -s /etc/vpn-gateway /var/lib/vpn-gateway/config
        fi

        echo ""
        echo "=========================================="
        echo "  VPN Gateway installed successfully!"
        echo "=========================================="
        echo ""
        echo "Next steps:"
        echo "  1. Run initial setup:"
        echo "     sudo vpn-gateway-setup [server-ip] [peer-count]"
        echo ""
        echo "  2. Start the service:"
        echo "     sudo systemctl start vpn-gateway"
        echo ""
        echo "  3. Enable on boot:"
        echo "     sudo systemctl enable vpn-gateway"
        echo ""
        ;;
esac

#DEBHELPER#
