# debian/vpn-gateway.service
[Unit]
Description=VPN Gateway - Transparent Proxy over WireGuard
Documentation=https://github.com/alexandergal86/vpn-gateway
After=network-online.target wireguard@wg0.service unbound.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/var/lib/vpn-gateway
Environment=RUST_LOG=vpn_gateway=info
Environment=CONFIG_PATH=/etc/vpn-gateway/gateway.json
ExecStartPre=/usr/lib/vpn-gateway/setup-network.sh start
ExecStart=/usr/bin/vpn-gateway
ExecStopPost=/usr/lib/vpn-gateway/setup-network.sh stop
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=no
ProtectSystem=strict
ReadWritePaths=/var/lib/vpn-gateway
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
